# defaults/profile.d/80-lliurex-desktop.sh
#!/bin/sh
export DCONF_PROFILE=lliurex
export XDG_CONFIG_DIRS="$XDG_CONFIG_DIRS:/etc/xdg/lliurex/desktop"
